# systemd/inverter-bridge.service
[Unit]
Description=Inverter Bridge - SRNE Modbus -> MQTT
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=inverter-bridge
Group=dialout
ExecStart=/opt/inverter-bridge/src/.venv/bin/python -m inverter_bridge --config /etc/inverter-bridge.yaml
Restart=always
RestartSec=5s
# The daemon pings WATCHDOG=1 once per hot cycle (~3s). 120s gives ample
# headroom for the worst-case wedged-serial recovery path (worker timeout).
WatchdogSec=120

# Security hardening (compatible with /dev/ttyUSB* + /etc + /var/log access)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/inverter-bridge /var/lib/inverter-bridge
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
